# crontab
PATH=/usr/local/bin:/usr/bin:/bin

# Personal combined: fetch battles + replays in one atomic pass
# Every 2 min — aggressive to ensure no replay data is lost
# Lockfile prevents overlap; personal data is highest priority
*/2 * * * * root /app/personal_combined.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Publish stats JSON to GitHub every 5 minutes
*/5 * * * * root /app/publish_wrapper.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Corpus battles: every minute, via the official CR API (IP-whitelisted,
# unaffected by RoyaleAPI/Cloudflare). High volume here is fine.
* * * * * root /app/corpus_scrape.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Corpus replays: SLOW trickle every 5 min (RoyaleAPI scrape through the VPN
# pool). Corpus-SCALE replay volume was proven non-viable 06-10 (rotation
# thrashed: 0 fetched / 30 challenges / 13 rotations in 6 min). This is the
# throttled opposite — a few players' freshest replays, gently — so the
# counterfactual sim gets OTHER players' games without spamming exits. Bumped
# 3→8 players / 10→5 min on 06-10 after a clean hour proved ample headroom.
2-59/5 * * * * root /app/corpus_replays.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: refresh top-ladder player tags (Monday 6am UTC)
0 6 * * 1 root /app/corpus_update.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Refresh simulation results daily at 2am UTC
0 2 * * * root /app/sim_refresh.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

# --- ML training ---
# Incremental WP inference: process new games every 5 min (offset from corpus)
*/5 * * * * root /app/wp_infer_new.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Incremental TCN embedding: embed new games every 6 hours
0 1,7,13,19 * * * root /app/embed_new.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: retrain TCN on full corpus (Sunday 4am UTC)
0 4 * * 0 root /app/tcn_train.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: retrain activity model for corpus scheduling (Monday 6:30am UTC, after corpus_update)
30 6 * * 1 root /app/train_activity.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

# --- Discovery jobs ---
# Hourly: add opponents I've lost to into the corpus
0 * * * * root /app/corpus_nemeses.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Daily: discover new players from opponent tags (3am UTC)
0 3 * * * root /app/corpus_discover.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: discover players from regional leaderboards (Monday 7am UTC)
0 7 * * 1 root /app/corpus_locations.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
