# crontab
# Fetch battles every minute for near-realtime dashboard updates
* * * * * /app/fetch.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Publish stats JSON to GitHub every 5 minutes
*/5 * * * * /app/publish_stats.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
