# crontab
PATH=/usr/local/bin:/usr/bin:/bin

# Personal combined: fetch battles + replays in one atomic pass
# Every 2 min — aggressive to ensure no replay data is lost
# Lockfile prevents overlap; personal data is highest priority
*/2 * * * * root /app/personal_combined.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Publish stats JSON to GitHub every 5 minutes
*/5 * * * * root /app/publish_wrapper.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Combined corpus scrape: battles + replays in one pass
# Fires every 5 min; lockfile ensures only one runs at a time
*/5 * * * * root /app/corpus_combined.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: refresh top-ladder player tags (Monday 6am UTC)
0 6 * * 1 root /app/corpus_update.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Refresh simulation results after corpus scrape (every 6h, offset 2h)
0 2,8,14,20 * * * root /app/sim_refresh.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

# --- ML training ---
# Incremental WP inference: process new games every 5 min (offset from corpus)
*/5 * * * * root /app/wp_infer_new.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: retrain TCN on full corpus (Sunday 4am UTC)
0 4 * * 0 root /app/tcn_train.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: retrain activity model for corpus scheduling (Monday 6:30am UTC, after corpus_update)
30 6 * * 1 root /app/train_activity.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

# --- Discovery jobs ---
# Hourly: add opponents I've lost to into the corpus
0 * * * * root /app/corpus_nemeses.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Daily: discover new players from opponent tags (3am UTC)
0 3 * * * root /app/corpus_discover.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: discover players from regional leaderboards (Monday 7am UTC)
0 7 * * 1 root /app/corpus_locations.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
