# crontab
PATH=/usr/local/bin:/usr/bin:/bin

# Fetch battles every minute for near-realtime dashboard updates
* * * * * root /app/fetch.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Publish stats JSON to GitHub every 5 minutes
*/5 * * * * root /app/publish_wrapper.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Fetch replay data from RoyaleAPI every 30 minutes
*/30 * * * * root /app/replay_wrapper.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Weekly: refresh top-ladder player tags (Monday 6am UTC)
0 6 * * 1 root /app/corpus_update.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Scrape battles for corpus players every 2 hours
0 */2 * * * root /app/corpus_scrape.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Refresh simulation results after corpus scrape (every 6h, offset 2h)
0 2,8,14,20 * * * root /app/sim_refresh.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Scrape replays for priority corpus players every 15 minutes
*/15 * * * * root /app/corpus_replays_priority.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
# Scrape replays for all corpus players every 30 minutes (30 players/run, rotates by staleness)
*/30 * * * * root /app/corpus_replays.sh >> /proc/1/fd/1 2>> /proc/1/fd/2
